# deploy/com.roger.metals-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.roger.metals-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/baga/Desktop/Roger/radar-potrebnostey/venv/bin/python</string>
        <string>/Users/baga/Desktop/Roger/radar-potrebnostey/bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/baga/Desktop/Roger/radar-potrebnostey</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/baga/Library/Logs/metals-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/baga/Library/Logs/metals-bot.log</string>
</dict>
</plist>
